# deployment/3sync.service
[Unit]
Description=3Sync - WebRTC Video Call App
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/3Sync
EnvironmentFile=/home/ubuntu/3Sync/.env
# Use source NVM and latest Node.js from PATH
ExecStart=/bin/bash -c 'source /home/ubuntu/.nvm/nvm.sh && exec node server.js'

# Auto-restart on failure
Restart=always
RestartSec=3

# Security
StandardOutput=journal
StandardError=journal
SyslogIdentifier=3sync

[Install]
WantedBy=multi-user.target
