# deployment/3sync.service
[Unit]
Description=3Sync - WebRTC Video Call App
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/3Sync
EnvironmentFile=/home/ubuntu/3Sync/.env
# Use node directly from PATH (installed via nvm as root)
ExecStart=/usr/local/bin/node server.js

# Auto-restart on failure
Restart=always
RestartSec=3

# Security
StandardOutput=journal
StandardError=journal
SyslogIdentifier=3sync

[Install]
WantedBy=multi-user.target
